# templates/aws-s3-sync.timer
[Unit]
Description=AWS Sync from S3 schedule timer

[Timer]
OnActiveSec=60s
OnUnitActiveSec=300s

[Install]
WantedBy=multi-user.target
